# build/pkg-debian/DEBIAN/postinst
#!/bin/sh
systemctl enable sender-smtp-api
systemctl status sender-smtp-api