# build/pkg-debian/DEBIAN/postinst
#!/bin/sh
systemctl enable project-identity
systemctl status project-identity